# etc/supervisor.ini
[inet_http_server]
port = 127.0.0.1:9001

[supervisord]
loglevel=info
directory=/var/log/supervisord
logfile=/var/log/supervisord/supervisord.log
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory  =  supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://127.0.0.1:9001

[program:gunicorn]
command=/usr/local/src/jdrf_env/bin/gunicorn -b 0.0.0.0:8000 wsgi:application 
directory=/usr/local/src/webservice/ccfa/web_cms
environment=VIRTUAL_ENV="/usr/local/src/jdrf_env",PATH="/usr/local/src/jdrf_env/bin:/usr/local/src/.linuxbrew/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/bin"
autostart=false
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisord/gunicorn.log

[program:nginx]
command = nginx 
user = root
autostart = false

